# templates/bredland/bredland-heartbeat.timer.template
[Unit]
Description=Run Bredland heartbeat every __BREDLAND_HEARTBEAT_INTERVAL__

[Timer]
OnBootSec=1min
OnUnitActiveSec=__BREDLAND_HEARTBEAT_INTERVAL__

[Install]
WantedBy=timers.target
