# templates/bredland/bredland-heartbeat.service.template
[Unit]
Description=Bredland telemetry heartbeat

[Service]
Type=oneshot
ExecStart=__BREDLAND_HEARTBEAT_SCRIPT_FILE__
